# automation/com.jmpark.sangfor.learnall.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd"\>
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jmpark.sangfor.learnall</string>

    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key><integer>2</integer>
      <key>Minute</key><integer>0</integer>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/jmpark/Documents/Playground/whelp99-code-sangfor-engineer-mcp/automation/scripts/run-learnall.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jmpark/Documents/Playground/whelp99-code-sangfor-engineer-mcp</string>

    <key>RunAtLoad</key>
    <false/>
  </dict>
</plist>
